# daily_call_data/crontab
# File: crontab
* * * * * python3 /app/main.py >> /app/cron.log 2>&1
# Dòng trống này rất quan trọng!